# layout/Library/LaunchDaemons/lh_moe.panties.discordextrasserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/DiscordExtrasServer</string>
	</array>
	<key>Label</key>
	<string>moe.panties.DiscordExtras</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>lh:moe.panties.discordextras</key>
		<true/>
	</dict>
</dict>
</plist>
