# layout/DEBIAN/postinst
#!/bin/sh

launchctl load /Library/LaunchDaemons/moe.panties.discordextrasserver.plist